# deploy/dsa110-operator-supervisor-h23.service
# The standing autonomy executor, running DIRECTLY ON h23 (no SSH tunnel).
#
# This is the recommended home for the single standing executor: it is
# headless (exposes NO port — nothing inbound to reach), and on h23 the etcd
# cluster and the dsa_monitor dashboard are local, so it needs no tunnel. It
# acquires the operator lease as session "supervisor" and runs the
# (policy-gated) health/recovery/injection/plan loops. Run AT MOST ONE of
# these across the whole site — it is the single executor.
#
# Install as a USER unit on h23 (run as the operator account, never root):
#   mkdir -p ~/.config/systemd/user
#   cp deploy/dsa110-operator-supervisor-h23.service ~/.config/systemd/user/
#   # put the env below in ~/.config/dsa110-operator/operator.env (chmod 600)
#   systemctl --user daemon-reload
#   systemctl --user enable --now dsa110-operator-supervisor-h23
#   loginctl enable-linger $USER     # keep it running after logout
#
# Required env (operator.env) for the on-h23, no-tunnel setup:
#   DSA_OPERATOR_ETCD_HOST=etcdv3service.pro.pvt   # reach etcd directly
#   DSA_OPERATOR_ETCD_PORT=2379
#   DSA_OPERATOR_DASHBOARD_PORT=5778               # dsa_monitor on h23 loopback
#   # Optional: DSA_OPERATOR_SLACK_WEBHOOK_URL=... for health alerts
#   # Optional: DSA_OPERATOR_ACTOR=supervisor       # audit label for its actions
#   # The supervisor is the deterministic, non-LLM loop: NO Anthropic key needed.
#   # Leave DSA_OPERATOR_ENFORCE_EGRESS unset (the host firewall is the control);
#   # if you do set it, the etcd host above is allowed automatically.
[Unit]
Description=dsa110-operator autonomy supervisor on h23 (single executor, no tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dsa110-operator
EnvironmentFile=%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
ExecStart=%h/dsa110-operator/.venv/bin/python -m dsa_operator.monitor.supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
