# deploy/dsa110-operator-supervisor-h23.service
# The standing autonomy executor, running DIRECTLY ON h23 (no SSH tunnel).
#
# This is the recommended home for the single standing executor: it is
# headless (exposes NO port — nothing inbound to reach), and on h23 the etcd
# cluster and the dsa_monitor dashboard are local, so it needs no tunnel. It
# acquires the operator lease as session "supervisor" and runs the
# (policy-gated) health/recovery/injection/plan loops. Run AT MOST ONE of
# these across the whole site — it is the single executor.
#
# IMPORTANT — paths are templated with %h (your home dir) and assume the repo
# is cloned at ~/dsa110-operator. If your clone lives elsewhere (e.g. on the
# DSA h23 box it is at ~/vikram/dev/dsa110-operator), edit WorkingDirectory,
# PYTHONPATH and ExecStart below to the ABSOLUTE path of your clone before
# installing. systemd will NOT expand shell variables here.
#
# Install as a USER unit on h23 (run as the operator account, never root):
#   mkdir -p ~/.config/systemd/user
#   cp deploy/dsa110-operator-supervisor-h23.service ~/.config/systemd/user/
#   # (optional) put extra env in ~/.config/dsa110-operator/operator.env (chmod 600)
#   systemctl --user daemon-reload
#   systemctl --user enable --now dsa110-operator-supervisor-h23
#   loginctl enable-linger $USER     # keep it running after logout
#
# The etcd/dashboard endpoints below are baked in as Environment= so the unit
# works with no operator.env at all. The optional EnvironmentFile is only for
# extras you may want later, e.g.:
#   # DSA_OPERATOR_SLACK_WEBHOOK_URL=...   # health alerts to Slack
#   # DSA_OPERATOR_ACTOR=supervisor        # audit label for its actions
# The supervisor is the deterministic, non-LLM loop: NO Anthropic key needed.
# Leave DSA_OPERATOR_ENFORCE_EGRESS unset (the host firewall is the control);
# if you do set it, the etcd host below is allowed automatically.
[Unit]
Description=dsa110-operator autonomy supervisor on h23 (single executor, no tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dsa110-operator
# Leading '-' makes this OPTIONAL: a missing file is NOT a startup failure.
EnvironmentFile=-%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/dsa110-operator/src
# etcd3 ships generated protobuf code older than its bundled protobuf runtime;
# force the pure-python impl so `import etcd3` doesn't raise at startup.
Environment=PROTOCOL_BUFFERS_PYTHON_IMPLEMENTATION=python
# Reach etcd + the dsa_monitor dashboard directly on h23 (no SSH tunnel).
Environment=DSA_OPERATOR_ETCD_HOST=etcdv3service.pro.pvt
Environment=DSA_OPERATOR_ETCD_PORT=2379
Environment=DSA_OPERATOR_DASHBOARD_PORT=5778
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
# PYTHONPATH=.../src above means no `pip install -e` is required — only the
# runtime deps (etcd3, pyyaml, requests) need to be in the chosen interpreter.
# ExecStart must point at a python that HAS those deps. Either:
#   (a) make a repo venv:  python3 -m venv .venv && .venv/bin/pip install -e '.[etcd]'
#       and keep the line below, or
#   (b) point it at an existing env that already has them, e.g. the dsart conda
#       env:  ExecStart=/opt/conda/envs/dsart_h23/bin/python -m dsa_operator.monitor.supervisor
ExecStart=%h/dsa110-operator/.venv/bin/python -m dsa_operator.monitor.supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
