# deploy/dsa110-operator-tunnel.service
# SSH tunnel to h23 (forwards etcd + dashboard to localhost).
#
# This is a TEMPLATE: @REPO@ and @PYTHON@ are filled in for your machine by
#   scripts/install_service.sh laptop
# which resolves your actual clone location and your conda/venv python. Do not
# `cp` this file directly unless you substitute those tokens yourself.
#
# Requires a passwordless SSH key for `h23` in ~/.ssh/config (Host h23 ...).
[Unit]
Description=dsa110-operator SSH tunnel to h23 (etcd + dashboard)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO@
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=@REPO@/src
ExecStart=@PYTHON@ -m dsa_operator.transport.ssh_tunnel --ssh-host h23
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
