# deploy/dsa110-operator-supervisor.service
# The standing autonomy executor. It acquires the operator lease as session
# "supervisor" and runs the (config-gated, off-by-default) health/recovery/
# injection/plan loops. Run AT MOST ONE of these across the whole site — it
# is the single executor. Leave loops disabled in config/policy.yaml until
# you intend the agent to act unprompted.
#
# This is a TEMPLATE: @REPO@ and @PYTHON@ are filled in for your machine by
#   scripts/install_service.sh laptop
# which resolves your actual clone location and your conda/venv python.
# (For the on-h23, no-tunnel deployment use the h23 mode of the installer.)
[Unit]
Description=dsa110-operator autonomy supervisor (single executor)
After=dsa110-operator-tunnel.service network-online.target
Wants=dsa110-operator-tunnel.service

[Service]
Type=simple
WorkingDirectory=@REPO@
# Leading '-' makes this OPTIONAL: a missing file is not a startup failure.
EnvironmentFile=-%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=@REPO@/src
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
ExecStart=@PYTHON@ -m dsa_operator.monitor.supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
