# deploy/dsa110-operator-web.service
# The web console (Flask + chat; local identity, no SSO). Depends on the tunnel.
#
# This is a TEMPLATE: @REPO@ and @PYTHON@ are filled in for your machine by
#   scripts/install_service.sh laptop
# which resolves your actual clone location and your conda/venv python.
#
# Secrets (ANTHROPIC_API_KEY, DSA_OPERATOR_SECRET_KEY, optional
# DSA_OPERATOR_SLACK_WEBHOOK_URL / DSA_OPERATOR_USER) come from an
# EnvironmentFile that is NOT in git (chmod 600). Set
# DSA_OPERATOR_ENFORCE_EGRESS=1 there to arm the in-process egress guard.
# Without the file the console still starts (read-only/stub agent); with no
# ANTHROPIC_API_KEY the chat falls back to the deterministic stub.
[Unit]
Description=dsa110-operator web console
After=dsa110-operator-tunnel.service network-online.target
Wants=dsa110-operator-tunnel.service

[Service]
Type=simple
WorkingDirectory=@REPO@
# Leading '-' makes this OPTIONAL: a missing file is not a startup failure.
EnvironmentFile=-%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=@REPO@/src
Environment=DSA_OPERATOR_BIND=127.0.0.1
Environment=DSA_OPERATOR_PORT=8787
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
ExecStart=@PYTHON@ -m dsa_operator.web.app
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
