# deploy/dsa110-operator-supervisor-h23.service
# The standing autonomy executor, running DIRECTLY ON h23 (no SSH tunnel).
#
# This is the recommended home for the single standing executor: it is
# headless (exposes NO port — nothing inbound to reach), and on h23 the etcd
# cluster and the dsa_monitor dashboard are local, so it needs no tunnel. It
# acquires the operator lease as session "supervisor" and runs the
# (policy-gated) health/recovery/injection/plan loops. Run AT MOST ONE of
# these across the whole site — it is the single executor.
#
# This is a TEMPLATE: @REPO@ and @PYTHON@ are filled in for your machine by
#   scripts/install_service.sh h23
# which resolves your actual clone location (anywhere — not just
# ~/dsa110-operator) and your conda/venv python. Do not `cp` this file
# directly unless you substitute those tokens yourself.
#
#   scripts/install_service.sh h23 --enable
#   loginctl enable-linger $USER     # keep it running after logout
#
# The etcd/dashboard endpoints below are baked in as Environment= so the unit
# works with no operator.env at all. The optional EnvironmentFile is only for
# extras you may want later, e.g.:
#   # DSA_OPERATOR_SLACK_WEBHOOK_URL=...   # health alerts to Slack
#   # DSA_OPERATOR_ACTOR=supervisor        # audit label for its actions
# The supervisor is the deterministic, non-LLM loop: NO Anthropic key needed.
# Leave DSA_OPERATOR_ENFORCE_EGRESS unset (the host firewall is the control);
# if you do set it, the etcd host below is allowed automatically.
[Unit]
Description=dsa110-operator autonomy supervisor on h23 (single executor, no tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO@
# Leading '-' makes this OPTIONAL: a missing file is NOT a startup failure.
EnvironmentFile=-%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=@REPO@/src
# etcd3 ships generated protobuf code older than its bundled protobuf runtime;
# force the pure-python impl so `import etcd3` doesn't raise at startup.
Environment=PROTOCOL_BUFFERS_PYTHON_IMPLEMENTATION=python
# Reach etcd + the dsa_monitor dashboard directly on h23 (no SSH tunnel).
Environment=DSA_OPERATOR_ETCD_HOST=etcdv3service.pro.pvt
Environment=DSA_OPERATOR_ETCD_PORT=2379
Environment=DSA_OPERATOR_DASHBOARD_PORT=5778
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
# @PYTHON@ must have the runtime deps (etcd3, pyyaml, requests). The installer
# defaults to your active conda/venv python; PYTHONPATH=@REPO@/src above means
# no `pip install -e` is needed — only those deps must be present.
ExecStart=@PYTHON@ -m dsa_operator.monitor.supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
