# deploy/dsa110-operator-supervisor.service
# The standing autonomy executor. It acquires the operator lease as session
# "supervisor" and runs the (config-gated, off-by-default) health/recovery/
# injection/plan loops. Run AT MOST ONE of these across the whole site — it
# is the single executor. Leave loops disabled in config/policy.yaml until
# you intend the agent to act unprompted.
[Unit]
Description=dsa110-operator autonomy supervisor (single executor)
After=dsa110-operator-tunnel.service network-online.target
Wants=dsa110-operator-tunnel.service

[Service]
Type=simple
WorkingDirectory=%h/dsa110-operator
EnvironmentFile=%h/.config/dsa110-operator/operator.env
Environment=PYTHONUNBUFFERED=1
Environment=DSA_OPERATOR_AUDIT_ROOT=%h/.local/share/dsa110-operator/audit
ExecStart=%h/dsa110-operator/.venv/bin/python -m dsa_operator.monitor.supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
